# examples/agenda-report/agenda-report.timer
[Unit]
Description=Run daily agenda-reports

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
